# scripts/systemd/k3s-metrics-push.service
[Unit]
Description=K3s cluster health metrics pusher (single-shot)
Documentation=https://github.com/ChadEngel/ce-ai-home-lab
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=cengel
Group=cengel
Environment=HOME=/home/cengel
WorkingDirectory=/home/cengel/source/ce-ai-home-lab
ExecStart=/home/cengel/source/ce-ai-home-lab/scripts/monitor_k3s_health.sh
TimeoutStartSec=90
# stdout/stderr -> journald (view with: journalctl -u k3s-metrics-push.service)

[Install]
WantedBy=multi-user.target